# slurm-resources/wasimoff_broker.service
[Unit]
Description=Starting the wasimoff broker
#After=
#Before=slurmctld.service

[Service]
ExecStart=WASIMOFF_ALLOWED_ORIGINS="*" WASIMOFF_HTTP_LISTEN=$(getent ahostsv4 controller | grep STREAM | sed -n '2p' | awk '{print $1}'):4080 go run ./
#ExecReload=
Type=simple
#Restart=

[Install]
WantedBy=default.target
#RequiredBy=slurmctld.service