# slurm-resources/wasimoff_broker.service
[Unit]
Description=Starting the wasimoff broker
#After=
Before=slurmctld.service

[Service]
# ExecStart=/bin/broker
ExecStart=/wasimoff_system/broker/broker
#ExecStop=
#ExecReload=
Type=exec
Restart=on-failure
Environment="WASIMOFF_ALLOWED_ORIGINS=\"*\""
Environment="WASIMOFF_HTTP_LISTEN=:4080"
KillMode=control-group
KillSignal=SIGTERM

[Install]
WantedBy=default.target
#RequiredBy=slurmctld.service