# slurm-resources/wasimoff_broker.service
[Unit]
Description=Starting the wasimoff broker
#After=
Before=slurmctld.service

[Service]
ExecStart=/usr/bin/bash /bin/broker/start_controller_node.sh
#ExecReload=
Type=simple
#Restart=

[Install]
WantedBy=default.target
RequiredBy=slurmctld.service